# pkg/deb/debian/preinst
#!/bin/bash
# maldet preinst — generated by pkg_gen.sh from pkg_lib templates
#
# Pre-install scriptlet: detects and backs up existing install.sh-based
# installations before package lays down FHS files.
#
# Placeholders:
#   maldet           — package name
#   /usr/local/maldetect   — legacy install path (e.g., /usr/local/bfd)
#   maldet       — binary name (e.g., bfd)
#
#
set -e

LEGACY_PATH="/usr/local/maldetect"

case "$1" in
    install|upgrade)
        # Detect any non-package state at $LEGACY_PATH that would conflict
        # with the symlink farm. The package ships these paths as symlinks;
        # if any exist as real (non-symlink) directories, dpkg will fail with
        # "trying to overwrite ... which is also in package" or similar
        # during file extraction. This catches install.sh installations,
        # partial installs, and hybrid states left by earlier package versions.
        _lmd_conflict=0
        if [ -f "$LEGACY_PATH/maldet" ] && [ ! -L "$LEGACY_PATH/maldet" ]; then
            _lmd_conflict=1
        else
            for _d in sigs quarantine sess tmp pub clean logs cron internals/alert; do
                if [ -d "$LEGACY_PATH/$_d" ] && [ ! -L "$LEGACY_PATH/$_d" ]; then
                    _lmd_conflict=1
                    break
                fi
            done
        fi
        if [ "$_lmd_conflict" -eq 1 ]; then
            _bkdir="${LEGACY_PATH}.bk.$(date +%Y%m%d-%s)"
            echo "Backing up existing installation to $_bkdir"
            cp -a "$LEGACY_PATH" "$_bkdir"
            rm -f "${LEGACY_PATH}.bk.last"
            ln -s "$_bkdir" "${LEGACY_PATH}.bk.last"
            # Preserve state files
            if [ -d "$LEGACY_PATH/tmp" ]; then
                mkdir -p /var/lib/maldet/tmp
                # cp may fail if tmp/ is empty — safe to ignore
                cp -a "$LEGACY_PATH"/tmp/* /var/lib/maldet/tmp/ 2>/dev/null || true
            fi
            # Stop SysVinit services if present — the old install.sh-based
            # installation may have registered init scripts that conflict with
            # the package's own service units; stop and deregister them so the
            # package post-install can set up services cleanly.
            for _initdir in /etc/rc.d/init.d /etc/init.d; do
                if [ -f "$_initdir/maldet" ]; then
                    # Service may not be running — safe to ignore stop failure
                    "$_initdir/maldet" stop 2>/dev/null || true
                    if command -v chkconfig >/dev/null 2>&1; then
                        # chkconfig may fail if service was never added — safe to ignore
                        chkconfig --del maldet 2>/dev/null || true
                    elif command -v update-rc.d >/dev/null 2>&1; then
                        # update-rc.d may fail if links already removed — safe to ignore
                        update-rc.d -f maldet remove 2>/dev/null || true
                    fi
                fi
            done
            # Stop systemd services before removing old install — service may
            # not exist or may not be running; either case is safe to ignore
            if command -v systemctl >/dev/null 2>&1; then
                systemctl stop maldet.service 2>/dev/null || true
            fi
            # Remove old install — package will lay down new files
            rm -rf "$LEGACY_PATH"
        fi
        # (no extra pre-install commands)
        ;;
esac

#DEBHELPER#

exit 0
